# scripts/homefinder.service
[Unit]
Description=HomeFinder Django Application
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/django-app
Environment=DJANGO_SETTINGS_MODULE=HomeFinderBackend.settings
ExecStart=/var/www/django-app/venv/bin/gunicorn --config /var/www/django-app/gunicorn_config.py HomeFinderBackend.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=2
StartLimitIntervalSec=10
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target